# servers/ci/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=java -jar jenkins.war
directory=/home/jenkins
environment = JENKINS_HOME="/var/jenkins_home"
